# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
user=root

[program:nginx]
command=nginx -g "daemon off;"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true

[program:page-stats-server]
command=/usr/local/bin/page-stats-server
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
environment=PORT="%(ENV_PAGE_STATS_PORT)s",REDIS_URL="%(ENV_REDIS_URL)s",APP_ENV="%(ENV_APP_ENV)s",HOST="%(ENV_PAGE_STATS_HOST)s"
